# gct-market-sentiment/gct-collector.service
[Unit]
Description=GCT Market Sentiment Continuous Collector
After=network.target

[Service]
Type=simple
User=$USER
WorkingDirectory=/Users/chris/GCT/gct-market-sentiment
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/Users/chris/.local/bin"
Environment="PYTHONPATH=/Users/chris/GCT/gct-market-sentiment"
ExecStart=/usr/bin/python3 /Users/chris/GCT/gct-market-sentiment/continuous_collector.py
Restart=always
RestartSec=10
StandardOutput=append:/Users/chris/GCT/gct-market-sentiment/logs/collector.log
StandardError=append:/Users/chris/GCT/gct-market-sentiment/logs/collector_error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target